# com.autoblog.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identifier -->
    <key>Label</key>
    <string>com.autoblog.daily</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/sethhondl/dev/active-projects/AutoBlog/scripts/login_trigger.py</string>
    </array>

    <!-- Run on login/load -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/sethhondl/dev/active-projects/AutoBlog</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/sethhondl/.nvm/versions/node/v22.9.0/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/sethhondl</string>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/autoblog.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/autoblog.err</string>

    <!-- Don't keep alive - just run once on load -->
    <key>KeepAlive</key>
    <false/>

    <!-- Nice level - run at low priority -->
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
